# maildir_idle.service
[Unit]
Description=IMAP IDLE agent

[Service]
ExecStart=maildir_idle_sync
KillSignal=SIGINT

[Install]
WantedBy=default.target
